# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:openldap]
command=/sbin/slapd -h 'ldap:/// ldapi:///' -g ldap -u ldap -d 0
